# resm-1/conf/debian/install.sh
#!/bin/bash
set -e

update-rc.d resm defaults
invoke-rc.d resm start
